# linux_things/crontab
*/30 * * * * rclone sync ~/panevka GoogleDrive:/panevka
